# pycharm-with-autoware/s6/services.d/projector/run
#!/usr/bin/with-contenv bash
source /home/autoware/Autoware/install/setup.bash
projector run