# legacy/pycharm-with-ros/s6/services.d/projector/run
#!/usr/bin/with-contenv bash
source /ros_entrypoint.sh
projector run