# pycharm/s6/services.d/projector/run
#!/usr/bin/with-contenv bash
projector run